# deploys/snapraid/templates/snapraid-scrub.timer.j2
[Unit]
Description=SnapRAID scrub timer

[Timer]
OnCalendar={{ schedule.scrub.day }} {{ schedule.scrub.time }}
Persistent=true

[Install]
WantedBy=timers.target
