# deploys/snapraid/templates/snapraid-sync.service.j2
[Unit]
Description=SnapRAID sync
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapraid -q sync
StandardOutput=journal
StandardError=journal
