# deploys/snapraid/templates/snapraid-smart.service.j2
[Unit]
Description=SnapRAID SMART check
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapraid -q smart
StandardOutput=journal
StandardError=journal
