# deploys/snapraid/templates/snapraid-smart.timer.j2
[Unit]
Description=SnapRAID SMART check timer

[Timer]
OnCalendar={{ schedule.smart.time }}
Persistent=true

[Install]
WantedBy=timers.target
