# deploys/snapraid/templates/snapraid-sync.timer.j2
[Unit]
Description=SnapRAID sync timer

[Timer]
OnCalendar={{ schedule.sync.time }}
Persistent=true

[Install]
WantedBy=timers.target
